# ansible/playbooks/templates/traefik/traefik-docker.service.j2
[Unit]
Description=Traefik via Docker Compose
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/etc/traefik
ExecStart=/usr/bin/docker compose -f /etc/traefik/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /etc/traefik/docker-compose.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target